# modules/install-etcd/etcd.path
[Path]
PathExists=/etc/sysconfig/etcd.conf

[Install]
WantedBy=multi-user.target