# modules/install-etcd/etcd-restart.path
[Path]
PathChanged=/opt/etcd/certs

[Install]
WantedBy=multi-user.target